# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/busybox-syslog.service
[Unit]
Description=System Logging Service

[Service]
Type=forking
RemainAfterExit=yes
ExecStart=/etc/initscripts/syslog start
ExecStop=/etc/initscripts/syslog stop

[Install]
WantedBy=multi-user.target
