# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/ipacmdiag.service
[Unit]
SourcePath=/usr/bin/ipacmdiag
Description= IPACM DIAG Service
Wants=diag-reboot-app.service chgrp-diag.service dev-diag.device
After=diag-reboot-app.service chgrp-diag.service dev-diag.device

[Service]
Restart=no
ExecStart=/usr/bin/ipacmdiag
User=radio
Group=radio
SupplementaryGroups=diag inet

[Install]
WantedBy=multi-user.target
