# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/mbimd.service
[Unit]
SourcePath=/usr/bin/mbimd
Description=QTI Mobile Broadband Service
After=mount-data.service

[Service]
Restart=always
ExecStart=/usr/bin/mbimd
User=radio
Group=radio
SupplementaryGroups=diag system kmsg wakelock

[Install]
WantedBy=multi-user.target