# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/mbim_fw_update@.service
[Unit]
Description= mbimd fw update helper service

[Service]
Environment="SCRIPT_ARGS=%I"
ExecStart=/etc/initscripts/mbim_fw_script.sh $SCRIPT_ARGS
Type=oneshot
RemainAfterExit=no