# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/mcm-core.service
[Unit]
Description= MCM RIL Services
SourcePath=/usr/bin/mcm_ril_service
Requires=postmss.service
After=postmss.service

[Service]
User=root
Group=root
Restart=no
ExecStartPre=/etc/initscripts/mcm_core_preload.sh
ExecStart=mcm_ril_service

[Install]
WantedBy=multi-user.target

