# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/qcmap_load_module@.service
[Unit]
Description= Load/Unload %I module service
SourcePath=%I
DefaultDependencies=no

[Service]
EnvironmentFile=/var/run/data/kernel_module.conf
Restart=no
Type=forking
RemainAfterExit=yes
ExecStart=insmod ${module}
ExecStop=rmmod ${module}
