# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/qmuxbridged.service
[Unit]
SourcePath=/usr/bin/qmuxbridge
Description=Qmuxbridge Service
After=mount-data.service

[Service]
Restart=always
ExecStart=/usr/bin/qmuxbridge
User=radio
Group=radio
SupplementaryGroups=diag system kmsg wakelock

[Install]
WantedBy=multi-user.target