# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/tinyproxyd.service
[Unit]
Description= tinyproxy service
SourcePath=/usr/sbin/tinyproxy


[Service]
Type=oneshot
RemainAfterExit=yes
Restart=no
ExecStart=/usr/sbin/tinyproxy -c /etc/data/tinyproxy.conf
ExecStop=/usr/bin/killall -15 tinyproxy
