# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/adbd.service.d/adbd.conf
[Unit]
Before=systemd-tmpfiles-setup-dev.service tmp.mount
