# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/usb.service.d/usb.conf
[Service]
Nice=-20

[Install]
WantedBy=local-fs.target
