# ql_crosstools/ql-ag520r-le17-gcc820-v1-toolchain/ql-sysroots/lib/systemd/system/usb.service.d/usb_data_service.conf
[Service]
ExecStartPre=/sbin/mount-copybind /systemrw/data/usb /etc/data/usb
